# release/jobs/nagios_dashboard/templates/nagios_reload.cron
*/5 * * * * sudo -u vcap /var/vcap/jobs/nagios_dashboard/bin/nagios_reload.sh
